# ops/systemd/omneuro-blog-sync.service
# ops/systemd/omneuro-blog-sync.service
[Unit]
Description=Omneuro Blog S3 → /var/www sync (oneshot)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/sysconfig/omneuro-blog
ExecStart=/usr/bin/env bash /home/ubuntu/omneuro/scripts/blog-sync.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target